# nginx_flask_python2/settings/supervisord.conf
[supervisord]
logfile = /dev/null
loglevel = info
pidfile = /var/run/supervisord.pid
nodaemon = true

[program:nginx]
command = nginx -g "daemon off;"

[program:uwsgi]
command = uwsgi --emperor /uwsgi/vassal/ --uid www-data --gid www-data

[program:setup]
command = python2 /settings/setup.py
startsecs = 0
autorestart = false
startretries = 1
